# BirdBath/beertaps/systemd/adc-reader-2.service
[Unit]
Description=ADC Reader for Controller 2 (0x49)
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/haven/BirdBath/beertaps
ExecStart=/home/pi/haven/BirdBath/beertaps/adc_reader.sh /home/pi/haven/BirdBath/beertaps/adc_config_2.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
